# modules/restic/autorestic.service
[Unit]
Description=Autorestic backup

[Service]
Type=oneshot
ExecStart=/usr/bin/autorestic -c /home/b/.config/autorestic/.autorestic.yml cron
